# qt4/context/etc/profile.d/qt4.sh
# Begin /etc/profile.d/qt4.sh

QT4DIR=/usr
export QT4DIR

# End /etc/profile.d/qt4.sh
